# bullseye.service
[Unit]
Description=Bulls Eye Predictive Engine (Streamlit)
After=network.target
OnFailure=bullseye-notify-failure.service

[Service]
Type=simple
User=YOUR_USERNAME
WorkingDirectory=YOUR_HOME_DIR/MSADS_/Time_series_31006/FINAL_PROJECT_VF

# Verify venv and key packages before starting
ExecStartPre=/bin/bash -c '\
    YOUR_HOME_DIR/MSADS_/Time_series_31006/tsenvi/bin/python -c \
    "import streamlit, pandas, numpy, matplotlib, xgboost" \
    || (echo "venv check failed" && exit 1)'

ExecStart=YOUR_HOME_DIR/MSADS_/Time_series_31006/tsenvi/bin/streamlit run app.py

Restart=on-failure
RestartSec=30
StartLimitIntervalSec=300
StartLimitBurst=3

StandardOutput=append:YOUR_HOME_DIR/MSADS_/Time_series_31006/FINAL_PROJECT_VF/logs/streamlit.log
StandardError=append:YOUR_HOME_DIR/MSADS_/Time_series_31006/FINAL_PROJECT_VF/logs/streamlit.log

[Install]
WantedBy=multi-user.target
